# server/ansible/neuron-catalog-meteor.conf
[program:neuron-catalog-meteor]
command=/usr/local/bin/meteor run --port 3450 --settings /neuron-catalog/server/config.json
directory=/neuron-catalog
